# daily-backup.timer
[Unit]
Description=Run daily-backup.service every day

[Timer]
OnCalendar=*.*.* 00:15:00