# daily-backup.service
[Unit]
Description=Run backup daily

[Service]
WorkingDirectory=/opt/bin/backup/
Type=simple
ExecStart=/usr/bin/bash backup-from.sh odoo-backup-daily
User=core

[Install]
WantedBy=multi-user.target
