# hourly-backup.timer
[Unit]
Description=Run hourly-backup.service every hour

[Timer]
OnCalendar=00,13:30:00
