# hourly-backup.timer
[Unit]
Description=Run hourly-backup.service every hour

[Timer]
OnCalendar=*:51:00
